# acme.sh/rootfs/etc/services.d/crond/run
#!/usr/bin/with-contenv bashio

# Run crond
bashio::log.info "Starting crond..."
exec crond -f -L /dev/stdout
